# serve.conf
description "My serve service"
author  "Sai<sai@foo.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /root/encoding
exec python serve.py
respawn

